# systemd/mihomo.service.d/override.conf
[Service]
# One-arm LAN gateway NTP bypass.
# Change this if the LAN/upstream interface is not eth1.
Environment=MIHOMO_LAN_IF=eth1
Environment=MIHOMO_NTP_RULE_PREF=8999

ExecStartPost=/bin/sh -c 'sysctl -w net.ipv4.ip_forward=1 >/dev/null'
ExecStartPost=/bin/sh -c 'ip rule add pref "$MIHOMO_NTP_RULE_PREF" ipproto udp dport 123 lookup main 2>/dev/null || true'
ExecStartPost=/bin/sh -c 'iptables -w -C FORWARD -i "$MIHOMO_LAN_IF" -o "$MIHOMO_LAN_IF" -p udp --dport 123 -j ACCEPT 2>/dev/null || iptables -w -I FORWARD 1 -i "$MIHOMO_LAN_IF" -o "$MIHOMO_LAN_IF" -p udp --dport 123 -j ACCEPT'
ExecStartPost=/bin/sh -c 'iptables -w -C FORWARD -i "$MIHOMO_LAN_IF" -o "$MIHOMO_LAN_IF" -m conntrack --ctstate ESTABLISHED,RELATED -j ACCEPT 2>/dev/null || iptables -w -I FORWARD 1 -i "$MIHOMO_LAN_IF" -o "$MIHOMO_LAN_IF" -m conntrack --ctstate ESTABLISHED,RELATED -j ACCEPT'
ExecStartPost=/bin/sh -c 'iptables -w -t nat -C POSTROUTING -o "$MIHOMO_LAN_IF" -p udp --dport 123 -j MASQUERADE 2>/dev/null || iptables -w -t nat -I POSTROUTING 1 -o "$MIHOMO_LAN_IF" -p udp --dport 123 -j MASQUERADE'
ExecStartPost=/bin/sh -c 'ip route flush cache || true'

ExecStopPost=/bin/sh -c 'ip rule del pref "$MIHOMO_NTP_RULE_PREF" 2>/dev/null || true'
ExecStopPost=/bin/sh -c 'iptables -w -D FORWARD -i "$MIHOMO_LAN_IF" -o "$MIHOMO_LAN_IF" -p udp --dport 123 -j ACCEPT 2>/dev/null || true'
ExecStopPost=/bin/sh -c 'iptables -w -D FORWARD -i "$MIHOMO_LAN_IF" -o "$MIHOMO_LAN_IF" -m conntrack --ctstate ESTABLISHED,RELATED -j ACCEPT 2>/dev/null || true'
ExecStopPost=/bin/sh -c 'iptables -w -t nat -D POSTROUTING -o "$MIHOMO_LAN_IF" -p udp --dport 123 -j MASQUERADE 2>/dev/null || true'
ExecStopPost=/bin/sh -c 'ip route flush cache || true'
